# deploy/systemd/maro-observe.service
[Unit]
Description=Maro Command Center — live dashboard + goal chat (port 80)
Documentation=file:///home/clawd/claude/maro-orchestration/CLAUDE.md
After=network.target

[Service]
Type=simple
User=clawd
WorkingDirectory=/home/clawd/claude/maro-orchestration
ExecStart=/usr/bin/python3 -m observe serve --port 80
Restart=on-failure
RestartSec=10

# Allow binding to port 80 without running as root
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

Environment=PYTHONPATH=/home/clawd/claude/maro-orchestration/src

StandardOutput=journal
StandardError=journal
SyslogIdentifier=maro-observe

[Install]
WantedBy=default.target
