# deploy/systemd/poe-observe.service
[Unit]
Description=Poe Command Center — live dashboard + goal chat (port 80)
Documentation=file:///home/clawd/claude/openclaw-orchestration/CLAUDE.md
After=network.target

[Service]
Type=simple
User=clawd
WorkingDirectory=/home/clawd/claude/openclaw-orchestration
ExecStart=/usr/bin/python3 -m observe serve --port 80
Restart=on-failure
RestartSec=10

# Allow binding to port 80 without running as root
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

Environment=PYTHONPATH=/home/clawd/claude/openclaw-orchestration/src
Environment=OPENCLAW_WORKSPACE=/home/clawd/.openclaw/workspace

StandardOutput=journal
StandardError=journal
SyslogIdentifier=poe-observe

[Install]
WantedBy=default.target
