# deploy/systemd/poe-heartbeat.service
[Unit]
Description=Poe Heartbeat — periodic health check + tiered self-healing
# Update Documentation path if repo location changes
Documentation=file:///home/clawd/claude/openclaw-orchestration/CLAUDE.md
After=network.target

[Service]
Type=simple
User=clawd
# NOTE: Update WorkingDirectory/PYTHONPATH if repo is moved.
# Use `scripts/heartbeat-ctl.sh` for path-independent startup instead.
WorkingDirectory=/home/clawd/claude/openclaw-orchestration
ExecStart=/usr/bin/python3 src/heartbeat.py --loop --interval 60
Restart=on-failure
RestartSec=30

# Token budget safety: stop after 4 hours, require manual restart.
# Prevents multi-day unattended token burn like the Apr 7-9 incident.
RuntimeMaxSec=14400

# Environment — PYTHONPATH must match WorkingDirectory
# No POE_MEMORY_DIR needed — defaults to ~/.poe/workspace/memory (same as config.py)
Environment=PYTHONPATH=/home/clawd/claude/openclaw-orchestration/src
Environment=OPENCLAW_WORKSPACE=/home/clawd/.openclaw/workspace

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=poe-heartbeat

[Install]
WantedBy=default.target
