# deploy/poe-telegram.service
[Unit]
Description=Poe Telegram Listener — routes Telegram messages through openclaw-orchestration handle()
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=clawd
WorkingDirectory=/home/clawd/claude/openclaw-orchestration/src
ExecStart=/usr/bin/python3 /home/clawd/claude/openclaw-orchestration/src/telegram_listener.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=poe-telegram

# Load credentials from openclaw.json (resolved at runtime by telegram_listener.py)
# Optionally override via env file:
# EnvironmentFile=/home/clawd/.openclaw/credentials/telegram.env

[Install]
WantedBy=multi-user.target
