# deploy/poe-inspector.service
[Unit]
Description=Poe Inspector — quality oversight loop
After=network.target

[Service]
Type=simple
User=clawd
WorkingDirectory=/home/clawd/claude/openclaw-orchestration/src
ExecStart=/usr/bin/python3 /home/clawd/claude/openclaw-orchestration/src/inspector.py --loop --interval 3600
Restart=on-failure
RestartSec=60
StandardOutput=journal
StandardError=journal
SyslogIdentifier=poe-inspector

[Install]
WantedBy=multi-user.target
